# container-openclaw.service
# container-openclaw.service
# OpenClaw bot — Mattermost gateway in a Podman container
#
# Installed by: install.sh (or manually copy to ~/.config/systemd/user/)
# Uses systemd specifiers: %h = home dir, %H = hostname — no editing needed.
# After changes: systemctl --user daemon-reload && systemctl --user restart container-openclaw.service
#
# Usage:
#   systemctl --user daemon-reload
#   systemctl --user enable --now container-openclaw.service
#   loginctl enable-linger $USER   # keep running after logout

[Unit]
Description=OpenClaw Bot (Podman container)
Documentation=man:podman-generate-systemd(1)
Wants=network-online.target
After=network-online.target
RequiresMountsFor=%t/containers

[Service]
Environment=PODMAN_SYSTEMD_UNIT=%n
Restart=always
RestartSec=10
TimeoutStopSec=70
ExecStart=/usr/bin/podman run \
	--cidfile=%t/%n.ctr-id \
	--cgroups=no-conmon \
	--rm \
	--sdnotify=conmon \
	--replace \
	-d \
	--name openclaw \
	--hostname %H \
	-e OPENCLAW_CONFIG_PATH=/config/openclaw.json \
	-v %h/openclaw-config:/config \
	-v %h/openclaw-workspace:/home/openclaw/workspace \
	-v %h/openclaw-state:/home/node/.openclaw \
	localhost/openclaw-jeeves:latest
ExecStop=/usr/bin/podman stop \
	--ignore -t 10 \
	--cidfile=%t/%n.ctr-id
ExecStopPost=/usr/bin/podman rm \
	-f \
	--ignore -t 10 \
	--cidfile=%t/%n.ctr-id
Type=notify
NotifyAccess=all

[Install]
WantedBy=default.target
